# rootfs/etc/init.d/rc.local
#!/bin/sh
echo "--------Start Local Services--------"


